# deploy/links.service
[Unit]
Description=Linksgood link shortener/lengthener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=infra
Group=infra
WorkingDirectory=/opt/links/app
Environment=NODE_ENV=production
EnvironmentFile=/etc/links/links.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
MemoryMax=220M

[Install]
WantedBy=multi-user.target
